# generated-configs/codersinflow-app.service
[Unit]
Description=Coders in Flow Application (Backend + Frontend)
After=network.target docker.service
Wants=docker.service

[Service]
Type=forking
User=root
WorkingDirectory=/var/www/codersinflow.com

# Environment variables
Environment="NODE_ENV=production"
Environment="PORT=4916"
Environment="HOST=0.0.0.0"
Environment="PUBLIC_API_URL=https://codersinflow.com"
Environment="BACKEND_PORT=8752"
Environment="MONGODB_URI=mongodb://admin:H+kmWOoxKC0yAOwaoimPyQ@127.0.0.1:27018/codersblog?authSource=admin"
Environment="JWT_SECRET=dv2ounuLU2BljGbtkSGg659GMpT21YNcmWthfSeSAHe5hSK4uB32XheEusuBMAI4auebg4g"
Environment="CORS_ORIGIN=https://codersinflow.com"
Environment="UPLOAD_DIR=/var/www/codersinflow.com/runtime/uploads"

# Start both services
ExecStartPre=/bin/bash -c 'docker compose -f /var/www/codersinflow.com/docker-compose.prod.yml up -d blog-mongodb'
ExecStartPre=/bin/bash -c 'sleep 5'
ExecStart=/bin/bash -c '/var/www/codersinflow.com/backend/server & echo $! > /var/www/codersinflow.com/runtime/backend.pid && /usr/bin/node /var/www/codersinflow.com/dist/server/entry.mjs & echo $! > /var/www/codersinflow.com/runtime/frontend.pid'

# Stop both services
ExecStop=/bin/bash -c 'kill $(cat /var/www/codersinflow.com/runtime/backend.pid) $(cat /var/www/codersinflow.com/runtime/frontend.pid) 2>/dev/null || true'
ExecStopPost=/bin/bash -c 'rm -f /var/www/codersinflow.com/runtime/backend.pid /var/www/codersinflow.com/runtime/frontend.pid'

Restart=always
RestartSec=10
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target