# generated-configs/codersinflow-backend.service
[Unit]
Description=Coders in Flow Backend API
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/codersinflow.com/backend
Environment="PORT=8752"
Environment="MONGODB_URI=mongodb://admin:H+kmWOoxKC0yAOwaoimPyQ@127.0.0.1:27421/codersblog?authSource=admin"
Environment="JWT_SECRET=dv2ounuLU2BljGbtkSGg659GMpT21YNcmWthfSeSAHe5hSK4uB32XheEusuBMAI4auebg4g"
Environment="CORS_ORIGIN=https://codersinflow.com"
Environment="UPLOAD_DIR=/var/www/codersinflow.com/runtime/uploads"
ExecStart=/var/www/codersinflow.com/backend/server
Restart=always
RestartSec=10
StandardOutput=append:/var/www/codersinflow.com/runtime/backend.log
StandardError=append:/var/www/codersinflow.com/runtime/backend.log

[Install]
WantedBy=multi-user.target